# templates/etc/profile.d/profile.j2
# {{ ansible_managed }}
#!/usr/bin/sh

{% if item.value.content is defined %}
{{ item.value.content }}
{% endif %}
